# packages/spirit-box_1.0.0_arm64/etc/systemd/system/disableSystemdLogging.service
[Unit]
Description=Disable systemd console logging.
StartLimitIntervalSec=0
#After=systemd
Before=system.slice dbus.service
#Before=console-setup.service

[Service]
Type=oneshot
RemainAfterExit=yes
# send signal to disable console logging
ExecStart=/bin/kill -s SIGRTMIN+21 1
# send signal to enable console logging
ExecStop=/bin/kill -s SIGRTMIN+20 1

[Install]
WantedBy=multi-user.target
