# packages/spirit-box_1.0.0_arm64/etc/systemd/system/spirit-box@ttys0.service
[Unit]
Description=TUI program for observing systemd units and script execution during Linux boot-up.
Conflicts=getty@%I.service
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/spirit-box
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/%I
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no
SendSIGHUP=yes

[Install]
WantedBy=multi-user.target
