# packages/spirit-box_1.0.0_arm64/etc/systemd/system/spirit-box.service
[Unit]
Description=Print systemd info to tty.
After=dbus.service disableSystemdLogging.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/spirit-box
StandardOutput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
