# packages/spirit-box_1.0.0_arm64/DEBIAN/postinst
#!/bin/sh

set -x

systemctl enable disableSystemdLogging.service
systemctl enable spirit-box@tty1.service
systemctl enable spirit-box@ttyS0.service
