# packages/spirit-box_1.0.0_arm64/etc/systemd/system/getty@tty1.service.d/override.conf
[Unit]
After=spirit-box.service

[Service]
TTYVTDisallocate=no
