# packages/spirit-box_1.0.0_arm64/etc/systemd/system/serial-getty@ttys0.service.d/override.conf
[Unit]
After=spirit-box@ttys0.service

[Service]
TTYVTDisallocate=no
